# docker/spark-standalone-cluster/conf/master.conf
[supervisord]
nodaemon=true
[program:spark-master]
command=/opt/spark/sbin/start-master.sh --properties-file /opt/conf/spark-defaults.conf
